# setup/systemd/tanya-web.service
[Unit]
Description=Tanya AI Streamlit Web Interface
After=network.target

[Service]
Type=simple
User=kalguy
WorkingDirectory=/home/kalguy/Tanya
ExecStart=/usr/bin/python3 -m streamlit run tanya_web.py --server.port 8501 --server.headless=true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
